# crontab
* * * * * npm run manage helloworld hi

# An empty line is required at the end of this file for a valid cron file.
